# dotfiles/git-annex/.config/autostart/git-annex-assistant.desktop
[Desktop Entry]
Type=Application
Name=git-annex-assistant
Exec=git annex assistant --autostart
StartupNotify=false
Terminal=false
Hidden=false
