# .crontab.sample
# set with `$ crontab -e'

@reboot screen -wipe
@reboot screen -dmS m
